# misc/myrpc-server/DEBIAN/prerm
systemctl disable --now myrpcd.service &
